# crontab
*/5 * * * * sudo bash /home/smarthome/test/solardataget.sh
